# dist/jnsk-web.service
[Unit]
Description=Jinsuke's website feat.Rocket

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/directory
Environment="ROCKET_ENV=prod"
Environment="ROCKET_ADDRESS=127.0.0.1"
Environment="ROCKET_PORT=8000"
Environment="ROCKET_LOG=critical"
ExecStart=/path/to/directory/target/release/two

[Install]
WantedBy=multi-user.target

